# pi-monitor-web.service
[Unit]
Description=Raspberry Pi Monitor Web Service
After=network.target

[Service]
Type=simple
User=pi
ExecStart=/usr/bin/python3 /usr/bin/pi_monitor_web.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
