# pi-monitor.service
[Unit]
Description=Raspberry Pi System Monitor
After=network.target

[Service]
Type=simple
User=peter
ExecStart=/usr/bin/python3 /usr/bin/pi_monitor.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
